# infra/deploy/com.agents.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agents.api</string>

  <!-- start_api.sh charge .env puis démarre uvicorn (remplace /path/to/project à l'install) -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/path/to/project/infra/deploy/start_api.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/project</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/path/to/project/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/path/to/project/logs/api.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/project/logs/api-error.log</string>

  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
